# server/globus_auth/systemd/globus-connector.service.template
[Unit]
Description=Globus Connector (ClearML token broker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@RUN_USER@
WorkingDirectory=@ROOT_DIR@
EnvironmentFile=@ROOT_DIR@/server/globus_auth/.env
Environment=PYTHONPATH=@ROOT_DIR@
Environment=PYTHONUNBUFFERED=1
ExecStart=@PYTHON_BIN@ -m uvicorn server.globus_auth.main:app --host 0.0.0.0 --port 8503
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
